# docker/cron/crontab
*/10 * * * * php /var/www/scripts/DeleteResetLinks.php >> /var/log/cron/cron.log 2>&1